# tee-broker/systemd/skill-library.service
[Unit]
Description=VerdantForged Skill Library Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/broker-daemon
EnvironmentFile=/opt/broker-daemon/config.env
ExecStart=/usr/local/bin/uvicorn skill_library.app:app --host 127.0.0.1 --port 8091 --workers 1
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/mnt/broker/skill-library /mnt/broker/logs

[Install]
WantedBy=multi-user.target
